# pkg/debian/DEBIAN/postinst
#!/bin/sh

set -x

systemctl daemon-reload
systemctl enable traefik
systemctl enable traefik-config
